# 00-ORCHESTRATION/launchd/com.syncrescendence.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/home/Desktop/syncrescendence/00-ORCHESTRATION/scripts/constellation_watchdog.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/home/Library/Logs/syncrescendence-watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/home/Library/Logs/syncrescendence-watchdog.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/home/.nvm/versions/node/v24.13.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/home</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/home/Desktop/syncrescendence</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
