# 00-ORCHESTRATION/scripts/launchd-mini/com.syncrescendence.watch-adjudicator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.watch-adjudicator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/home/Desktop/syncrescendence/00-ORCHESTRATION/scripts/watch_dispatch.sh</string>
        <string>adjudicator</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/home/Desktop/syncrescendence</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/syncrescendence-watch-adjudicator.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syncrescendence-watch-adjudicator.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/home</string>
        <key>NODE_NO_WARNINGS</key>
        <string>1</string>
        <key>SYNCRESCENDENCE_CODEX_MODEL</key>
        <string>gpt-5.2-codex</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
