# corpus/11443.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.circadian-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/system/syncrescendence/orchestration/00-ORCHESTRATION/scripts/circadian_sync.py</string>
        <string>--repo-root</string>
        <string>/Users/system/syncrescendence</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/system/syncrescendence</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SYNCRESCENDENCE_ROOT</key>
        <string>/Users/system/syncrescendence</string>
        <key>HOME</key>
        <string>/Users/system</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/system/syncrescendence/orchestration/state/circadian.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/system/syncrescendence/orchestration/state/circadian.stderr.log</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
