# corpus/11669.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.ingest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$(git -C ~/syncrescendence rev-parse --show-toplevel)" && ./scripts/auto_ingest_loop.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/system/syncrescendence/memory/ingest-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/system/syncrescendence/memory/ingest-stderr.log</string>
</dict>
</plist>
