# corpus/11671.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/system</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_ADJUDICATOR</key>
		<string>mini</string>
		<key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_AJNA</key>
		<string>local</string>
		<key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_CARTOGRAPHER</key>
		<string>mini</string>
		<key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_COMMANDER</key>
		<string>mini</string>
		<key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_PSYCHE</key>
		<string>mini</string>
	</dict>
	<key>Label</key>
	<string>com.syncrescendence.proactive-orchestrator</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/system/syncrescendence/orchestration/00-ORCHESTRATION/scripts/proactive_orchestrator.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/system/Library/Logs/syncrescendence-orchestrator.log</string>
	<key>StandardOutPath</key>
	<string>/Users/system/Library/Logs/syncrescendence-orchestrator.log</string>
	<key>StartInterval</key>
	<integer>307</integer>
	<key>WorkingDirectory</key>
	<string>/Users/system/syncrescendence</string>
</dict>
</plist>
