# com.syncrescendence.ledger-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.ledger-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/system/syncrescendence/orchestration/00-ORCHESTRATION/scripts/ledger_refresh.sh</string>
        <string>--report-md</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SYNCRESCENDENCE_REPO_ROOT</key>
        <string>/Users/system/syncrescendence</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>8</integer>
        <key>Minute</key><integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/system/syncrescendence/orchestration/state/ledger-refresh-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/system/syncrescendence/orchestration/state/ledger-refresh-stderr.log</string>
</dict>
</plist>
