# infra/SCRIPT-ORCHESTRATION-LAUNCHD-com.syncrescendence.memsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.memsync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/system/syncrescendence/orchestration/scripts/memsync_daemon.py</string>
        <string>--repo-root</string>
        <string>/Users/system/syncrescendence</string>
        <string>--graphiti-base</string>
        <string>http://M1-Mac-mini.local:8001</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/system/syncrescendence</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/system/syncrescendence/orchestration/state/memsync.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/system/syncrescendence/orchestration/state/memsync.stderr.log</string>
</dict>
</plist>
