# corpus/design-taste/11674.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.youtube-ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/system/.syncrescendence/venv/bin/python3</string>
        <string>/Users/system/syncrescendence/orchestration/00-ORCHESTRATION/scripts/youtube_ingest.py</string>
        <string>--ingest-queue</string>
        <string>--rebuild-index</string>
    </array>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/Users/system/.syncrescendence/logs/youtube_ingest.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/system/.syncrescendence/logs/youtube_ingest.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/system</string>
        <key>SYNCRESCENDENCE_CONFIG_DIR</key>
        <string>/Users/system/.syncrescendence</string>
        <key>SYNCRESCENDENCE_SOURCES_DIR</key>
        <string>/Users/system/syncrescendence/sources/04-SOURCES</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
