# 00-ORCHESTRATION/scripts/launchd/com.syncrescendence.sensing-frontier-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.sensing-frontier-scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/home/.syncrescendence/scripts/run_claude_task.sh</string>
        <string>frontier-scan</string>
        <string>/Users/home/.syncrescendence/claude-tasks/frontier-scan.md</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/home/Desktop/syncrescendence</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/syncrescendence-claude/sensing-frontier-scan.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syncrescendence-claude/sensing-frontier-scan.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/home</string>
        <key>PATH</key>
        <string>/Users/home/.nvm/versions/node/v24.13.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
