# 00-ORCHESTRATION/scripts/launchd/com.syncrescendence.claude-linear-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.claude-linear-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/home/.syncrescendence/scripts/run_claude_task.sh</string>
        <string>linear-status-check</string>
        <string>/Users/home/.syncrescendence/claude-tasks/linear-status-check.md</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/home/Desktop/syncrescendence</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>7</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>19</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/syncrescendence-claude/claude-linear-check.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syncrescendence-claude/claude-linear-check.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/home</string>
        <key>PATH</key>
        <string>/Users/home/.nvm/versions/node/v24.13.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
