# 00-ORCHESTRATION/scripts/launchd/com.syncrescendence.watch-canon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.watch-canon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/home/Desktop/syncrescendence/00-ORCHESTRATION/scripts/watch_canon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/home/Desktop/syncrescendence</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/syncrescendence-watch-canon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syncrescendence-watch-canon.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_NO_WARNINGS</key>
        <string>1</string>
        <key>NODE_OPTIONS</key>
        <string>--no-warnings</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
