# scaffold/scripts/orchestration/com.syncrescendence.dag-tension-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.dag-tension-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/system/syncrescendence/orchestration/00-ORCHESTRATION/scripts/dag_tension_monitor.py</string>
        <string>--repo-root</string>
        <string>/Users/system/syncrescendence</string>
        <string>--mode</string>
        <string>monitor</string>
    </array>

    <!-- Run every 6 hours (21600 seconds) per Adjudicator spec -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/Users/system/syncrescendence/orchestration/00-ORCHESTRATION/state/DYN-DAG_TENSION_MONITOR.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/system/syncrescendence/orchestration/00-ORCHESTRATION/state/DYN-DAG_TENSION_MONITOR.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/system/syncrescendence</string>

    <!-- CRITICAL: launchd does NOT source ~/.zshrc. Set env vars here. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/system/syncrescendence</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
