# corpus/infrastructure/09107_from_infrastructure.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.claude-corpus-insight</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/home/.syncrescendence/scripts/run_claude_task.sh</string>
        <string>corpus-insight</string>
        <string>/Users/home/.syncrescendence/claude-tasks/corpus-insight.md</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/home/Desktop/syncrescendence</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/syncrescendence-claude/claude-corpus-insight.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syncrescendence-claude/claude-corpus-insight.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/home</string>
        <key>PATH</key>
        <string>/Users/home/.nvm/versions/node/v24.13.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
