# scaffold/scripts/orchestration/launchd-mini/com.syncrescendence.auto-ingest-supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syncrescendence.auto-ingest-supervisor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/home/Desktop/syncrescendence/orchestration/scripts/auto_ingest_supervisor.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/home/Desktop/syncrescendence</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/home</string>
        <key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_AJNA</key>
        <string>macbook-air</string>
        <key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_COMMANDER</key>
        <string>local</string>
        <key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_ADJUDICATOR</key>
        <string>local</string>
        <key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_CARTOGRAPHER</key>
        <string>local</string>
        <key>SYNCRESCENDENCE_REMOTE_AGENT_HOST_PSYCHE</key>
        <string>local</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/syncrescendence-auto-ingest-supervisor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syncrescendence-auto-ingest-supervisor.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
